# deploy/alpaca-crypto.timer
[Unit]
Description=Run Alpaca Crypto Strategy Engine every 4 hours

[Timer]
# Crypto runs 24/7, every day, every 4 hours
# Unlike equities, no weekday restriction — crypto never sleeps
OnCalendar=*-*-* 00/4:00:00
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
